# usr/lib/systemd/system/train-tickets-frontend.service
[Unit]
Description=Train Tickets Frontend (Vite Dev Server)
After=network.target

[Service]
Type=simple
WorkingDirectory=/usr/share/train-tickets/frontend
ExecStart=/usr/bin/npx vite --host 0.0.0.0
Restart=on-failure
RestartSec=5

EnvironmentFile=-/etc/train-tickets/frontend.env

[Install]
WantedBy=multi-user.target
